# systemd/kdx-identity.service
[Unit]
Description=NAR Object Storage appliance identity initializer
After=local-fs.target
Before=kdx-agent.service kdx-wizard.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/kronosdx/scripts/initialize-appliance-identity.sh

[Install]
WantedBy=multi-user.target
